# ModeloMACOS/launchd/ollama.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.piper.ollama</string>
  <key>ProgramArguments</key>
  <array>
    <string>ollama</string>
    <string>serve</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OLLAMA_HOST</key><string>127.0.0.1:11434</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>${HOME}/.local/share/piper-cli/logs/ollama.out.log</string>
  <key>StandardErrorPath</key><string>${HOME}/.local/share/piper-cli/logs/ollama.err.log</string>
</dict>
</plist>
